# systemd/venator-c2.service
[Unit]
Description=Venator C2 server — LoRa command link to the drone
# Do NOT wait for multi-user.target: on this Jetson, snapd.seeded.service stalls
# it for ~2 min at boot. We only need the USB ports (which come up early) and we
# retry until the Pixhawk is ready, so start as soon as the filesystem is up.
After=local-fs.target
# Never give up restarting while it waits for the ports / Pixhawk to appear.
StartLimitIntervalSec=0

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/Desktop/pixhawk-test
# Real flight controller, props OFF by default: only motor tests are runnable
# remotely; flight missions stay rejected until this is deliberately changed to
# add --props-on. Nothing arms on boot — arming only happens on a confirmed
# command that also passes PX4's own pre-arm checks.
ExecStart=/home/jetson/Desktop/pixhawk-test/venv/bin/python /home/jetson/Desktop/pixhawk-test/jetson_c2_server.py --real
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
